# templates/systemd.timer.j2
[Unit]
Description=Timer for {{ item.name }}
Requires=restic-{{ item.name }}.service

[Timer]
Unit=restic-{{ item.name }}.service
OnCalendar={{ item.calendar }}

[Install]
WantedBy=timers.target
